# systemd/pu1des-coletor.service
[Unit]
Description=PU1DES Status Collector

[Service]
Type=oneshot
ExecStart=/usr/bin/php /opt/pu1des-reflector/scripts/coletar_status.php
